# com.statsports.test-extract.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.statsports.test-extract</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/christopherdickson/Documents/GitHub/statsports_starter/cron_test.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>WorkingDirectory</key>
    <string>/Users/christopherdickson/Documents/GitHub/statsports_starter</string>
    <key>StandardOutPath</key>
    <string>/tmp/statsports_test_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/statsports_test_launchd_error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>